# deploy/stocksignal-frontend.service
[Unit]
Description=StockSignal Next.js Frontend
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/stocksignal/frontend
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node /opt/stocksignal/frontend/node_modules/.bin/next start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
